# systemd/cs2bot.service
[Unit]
Description=CS2 Discord manager bot (owns and supervises the CS2 server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=steam
WorkingDirectory=/home/steam/cs2bot
Environment=CS2BOT_CONFIG=/home/steam/cs2bot/config.json
ExecStart=/home/steam/cs2bot/venv/bin/python -m cs2bot.bot
# The bot launches CS2 as a child process, so give it time to stop cleanly.
TimeoutStopSec=60
KillMode=mixed
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
